# systemd/ride-dispatch.service
[Unit]
Description=Production Lab Service C — Ride Dispatch Service (port 3003)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
Group=YOUR_USERNAME
WorkingDirectory=YOUR_PROJECT_PATH/services/service-c

# All config via environment — nothing hardcoded in the app
Environment=BIND_HOST=127.0.0.1
Environment=SERVICE_PORT=3003
Environment=CALLBACK_URL=http://service-a.internal:3001/ride/callback
Environment=DOWNSTREAM_TIMEOUT=5
Environment=PYTHONPATH=YOUR_PROJECT_PATH/services

ExecStart=YOUR_PROJECT_PATH/venv/bin/uvicorn app:app --host 127.0.0.1 --port 3003

Restart=on-failure
RestartSec=2

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
